# system_configuration/profile.d/bluetooth.sh
#!/bin/bash

/usr/lib/bluetooth/bluetooth_restart &
