# system_configuration/profile.d/eina-log.sh
export EINA_LOG_LEVEL=0
export EINA_LOG_LEVELS=e:0
export EINA_LOG_BACKTRACE=0
export EINA_LOG_LEVEL_MAXIMUM=0
