# system_configuration/profile.d/fcitx.sh
# Fcitx Setting
#export LANG="zh_CN.UTF-8" # Setting Chinese locales.
#export XIM_PROGRAM=fcitx
#export XIM=fcitx
#export GTK_IM_MODULE=fcitx
#export QT_IM_MODULE=fcitx
#export XMODIFIERS=@im=fcitx

GTK_IM_MODULE=fcitx
QT_IM_MODULE=fcitx
XMODIFIERS=@im=fcitx
INPUT_METHOD=fcitx
SDL_IM_MODULE=fcitx
GLFW_IM_MODULE=ibus
