# system_configuration/profile.d/fcitx.sh
# Setting Chinese locales.
#export LANG="zh_CN.UTF-8"

# Fcitx Setting
##export GLFW_IM_MODULE=ibus
#export GLFW_IM_MODULE=fcitx
#export XMODIFIERS=@im=fcitx
#export GTK_IM_MODULE=fcitx
#export SDL_IM_MODULE=fcitx
#export INPUT_METHOD=fcitx
#export QT_IM_MODULE=fcitx
#export XIM_PROGRAM=fcitx
#export XIM=fcitx

#export GLFW_IM_MODULE=ibus
GLFW_IM_MODULE=fcitx
XMODIFIERS=@im=fcitx
GTK_IM_MODULE=fcitx
SDL_IM_MODULE=fcitx
INPUT_METHOD=fcitx
QT_IM_MODULE=fcitx
XIM_PROGRAM=fcitx
XIM=fcitx
