# system_configuration/profile.d/ibus.sh
# iBus always on 
# create home/user/.xprofile, add following commands
#export GTK_IM_MODULE=ibus
#export XMODIFIERS=@im=ibus
#export QT_IM_MODULE=ibus
#ibus-daemon -drx

#GTK_IM_MODULE=ibus
#XMODIFIERS=@im=ibus
#QT_IM_MODULE=ibus
#ibus-daemon -drx
