# system_configuration/profile.d/kodi
# This is /etc/profile.d/kodi        ETP 21st August 2016
# Fix to launch Kodi
export KODI_HOME=/usr/share/kodi/