# system_configuration/profile.d/pigpio_default_address.sh
export GPIOZERO_PIN_FACTORY=pigpio
export PIGPIO_ADDR=fe80::1%usb0
